# vagrant/cookbook/templates/default/jackrabbit.erb
#!/bin/bash
### BEGIN INIT INFO
# Provides:          jackrabbit
# Short-Description: Start/stop Jackrabbit JCR server.
#
# Description:      This relies on a PID file to check if Jackrabbit is running.
#                   If you kill Jackrabbit without removing the PID file, you
#                   will not be able to start Jackrabbit with this script until
#                   you manually remove the PID file.
#                   Edit the variables below to configure Jackrabbit
#                   Depending on the storage backend, you might want to adjust
#                   the required start / stop lines.
#
# Default-Start:    2 3 4 5
# Default-Stop:     0 1 6
# Required-Start:
# Required-Stop:
#
# Author:           Daniel Barsotti <daniel.barsotti@liip.ch>
#                   Bastian Widmer <bastian.widmer@liip.ch>
#
#                   Modified by <waldvogel@liip.ch>
### END INIT INFO


### PLEASE EDIT THESE VALUES FOR YOUR SETUP
# Directory in which jackrabbit will store local data, must be writeable
BASEDIR=/opt/jackrabbit
# Full filename of jackrabbit standalone .jar to run
JACKRABBIT_JAR=$BASEDIR/jackrabbit.jar
# IP address for jackrabbit to listen on. you can make jackrabbit listen on all
# interfaces by using 0.0.0.0 here.
JACKRABBIT_HOST=0.0.0.0
# Port number to listen on. 8080 is default, you can use something else
JACKRABBIT_PORT=8080
# Java memory allocation
MEMORY="-XX:MaxPermSize=128m -Xmx512M -Xms128M"
# Full filename to the logfile to output console output of jackrabbit
# additionally, jackrabbit has its own logfile
LOGFILE=$BASEDIR/jackrabbit.log
###

PID=$(pgrep -f "java $MEMORY -jar $JACKRABBIT_JAR -h $JACKRABBIT_HOST -p $JACKRABBIT_PORT")

do_start() {
    if [ -z "$PID" ] ; then
        cd $BASEDIR
        nohup java $MEMORY -jar $JACKRABBIT_JAR -h $JACKRABBIT_HOST -p $JACKRABBIT_PORT >> $LOGFILE 2>&1 &
        # Wait until the server is ready (from an idea of Christoph Luehr)
        echo "Jackrabbit started in the background"
    else
        echo "Jackrabbit is already running"
    fi
}

do_stop() {
    if ! [ -z "$PID" ] ; then
        kill "$PID"
        echo "Jackrabbit stopped"
    else
        echo "Jackrabbit is not running"
    fi
    exit 3
}

do_status() {
    if ! [ -z "$PID" ]; then
          echo "Jackrabbit is running [ pid = "$PID"]"
    else
        echo "Jackrabbit is not running"
        exit 3
    fi
}

case "$1" in
  start)
        do_start
        ;;
  stop)
        do_stop
        ;;
    status)
        do_status
        ;;
  *)
        echo "Usage: $SCRIPTNAME {start|stop|status}" >&2
        exit 3
    ;;
esac
